# common/ansible/roles/system/config/files/logrotate.timer.d/override.conf
[Timer]
OnCalendar=
OnCalendar=*-*-* 0,12:00:00
